# deployment/orgos-gpu.service
# =============================================================================
# deployment/orgos-gpu.service
# systemd service file for the FastAPI inference gateway.
# This serves ALL Dragnet AI systems — OrgOS, HR, Finance, and others.
#
# Install on RunPod pod:
#   cp deployment/orgos-gpu.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable orgos-gpu
#   systemctl start orgos-gpu
#
# Useful commands:
#   systemctl status orgos-gpu       — check if running
#   journalctl -u orgos-gpu -f       — follow live logs
#   systemctl restart orgos-gpu      — restart after a code update
#   systemctl stop orgos-gpu         — stop before shutting down the pod
# =============================================================================

[Unit]
Description=Dragnet GPU Inference Gateway (FastAPI) — serves all Dragnet AI systems
Documentation=https://github.com/Dragnet-Solutions/orgos-gpu
After=network.target redis.service

[Service]
User=dragnet-deploy
Group=dragnet-deploy

WorkingDirectory=/workspace/orgos-gpu

# Load environment variables from the .env file on the pod.
# This file is NEVER committed to GitHub — it contains secrets.
EnvironmentFile=/workspace/orgos-gpu/.env

# Single uvicorn worker — Celery handles all concurrency
ExecStart=/usr/local/bin/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 1 \
    --log-level info

Restart=on-failure
RestartSec=5

# Allow 120 seconds for startup before systemd marks it as failed
TimeoutStartSec=120

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target